# docker/crontab
* * * * * root supervisorctl restart workflow
*/15 * * * * root supervisorctl restart stats
