# docker/supervisor/workflow.conf
[supervisord]

[program:workflow]
user=cerberus
command=/usr/bin/python /home/cerberus/cerberus-core/event/workflow.py
process_name=%(program_name)s
autostart=true
autorestart=false
startretries=0
numprocs=1
priority=100
stderr_logfile=/home/cerberus/cerberus-core/workflow.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
